# systemd/neflare-cn-ssh-geo-update.service
[Unit]
Description=Refresh mainland-China SSH geo-block nftables sets
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neflarectl update-cn-ssh-geo

